# board/common/rootfs-overlay/etc/inittab
# Documentación acerca del formato de este fichero disponible en:
# https://manpages.debian.org/testing/sysvinit-core/inittab.5.en.html

id:2:initdefault:

mnt::sysinit:/bin/mount -a
swp::sysinit:/sbin/swapon -a
rcS:2:wait:/etc/init.d/rcS

# Permitir que root se conecte desde el puerto serie
S0:12:respawn:/sbin/sulogin -p -e

# Comandos a ejecutar cuando el sistema se apaga
rcK:06:wait:/etc/init.d/rcK
swpo:06:wait:/sbin/swapoff -a
# Al reiniciar el sistema, kexec-reboot se encargará
# de desmontar los sistemas de ficheros
mnto:0:wait:/bin/umount -ar

# Comandos a ejecutar para manejar apagado y reinicio
hlt:0:wait:/sbin/poweroff -h
reb:6:wait:/sbin/kexec-reboot
