# board/pinetainer/rootfs-overlay/etc/network/if-pre-up.d/wlan
#!/bin/sh

# Iniciar wpa_supplicant para configurar el enlace inalámbrico
case $IFACE in
	wlan*)
		wpa_supplicant -i$IFACE -B -c/etc/wpa_supplicant.conf -s || exit $?;;
esac
