# board/common/rootfs-overlay/etc/profile.d/gzip-aliases.sh
#!/bin/sh

# Use pigz as a drop-in replacement for gzip
alias gzip='pigz'
alias gunzip='pigz -d'
