# board/common/rootfs-overlay/etc/profile.d/pigz.sh
alias gzip='pigz'
alias gunzip='pigz -d'
